# distro/systemd/rg-status-bar.service
# rg-status-bar.service — Refresh fleet status cache for i3blocks
#
# Runs rg-status-bar.sh --all --json and writes to /tmp/rg-status.json.
# i3blocks reads from this cache file instead of invoking the binary directly,
# preventing bar freezes during slow MCP calls.
#
# Paired with rg-status-bar.timer (30s interval).

[Unit]
Description=Refresh ralphglasses fleet status cache for status bar
After=network-online.target ralphglasses.service
Wants=network-online.target

[Service]
Type=oneshot
User=ralph
ExecStart=/bin/sh -c '/opt/ralphglasses/distro/scripts/rg-status-bar.sh --all --json > /tmp/rg-status.json 2>/dev/null || true'
Restart=on-failure
RestartSec=5
MemoryMax=64M
CPUQuota=10%
