# distro/systemd/watchdog.service
[Unit]
Description=ralphglasses Watchdog
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ralphglasses-watchdog
Restart=always
RestartSec=5
WatchdogSec=30
NotifyAccess=all
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000
User=ralph

[Install]
WantedBy=graphical-session.target
